# deploy_tools/gunicorn-systemd.template.service
[Unit]
Description=Gunicorn server for SITENAME

[Service]
Restart=on-failure
User=ubuntu